# scheduled/tidy_radio_directory.cron
#!/bin/bash

# Say when we want to remove files from Radio directory
/usr/bin/find /home/pi/Radio/tmp -mtime +10 -exec rm -f {} \;
