# scheduled_tasks/pi.cron
# retrieve programs
# M H DoM M DoW
3 0 * * * /home/pi/smartbot/get_iplayer/get_radio_shows.sh > /home/pi/smartbot/log/get_iplayer.log 2>&1
13 0 * * * /home/pi/smartbot/bin/download_facts.sh > /home/pi/smartbot/log/download_facts.log 2>&1
23 0 * * * /home/pi/smartbot/scheduled_tasks/tidy_radio_directory.cron > /home/pi/smartbot/log/tidy_radio_directory.log 2>&1
03 1 * * * /usr/bin/python3 /home/pi/smartbot/bin/getpodcasts.py > /home/pi/smartbot/log/getpodcasts.log 2>&1
# Update raspian (Saturday early morning!)
13 3 * * 6 /home/pi/smartbot/bin/updatepi.sh

