# scheduled_tasks/pi.cron
# retrieve programs
1 2 * * * /home/pi/smartbot/get_iplayer/get_radio_shows.sh > /home/pi/smartbot/log/get_iplayer.log 2>&1
1 10 * * * /home/pi/smartbot/bin/download_facts.sh > /home/pi/smartbot/log/download_facts.log 2>&1