# scheduled_tasks/updatepi.cron
# Update raspian
1 3 * * 6 /home/pi/smartbot/bin/updatepi.sh

